# ops/launchd/com.poptonic.ocr-insurance-knowledge-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.poptonic.ocr-insurance-knowledge-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/wenshuping/Documents/OCR_insurance/scripts/run-insurance-knowledge-daily.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/wenshuping/Documents/OCR_insurance</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/wenshuping/Documents/OCR_insurance/.runtime/launchd-insurance-refresh.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/wenshuping/Documents/OCR_insurance/.runtime/launchd-insurance-refresh.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>SCRAPLING_PYTHON_BIN</key>
    <string>/Users/wenshuping/Documents/Scrapling/.venv/bin/python</string>
    <key>SCRAPLING_PROJECT_DIR</key>
    <string>/Users/wenshuping/Documents/Scrapling</string>
  </dict>
</dict>
</plist>
